# deploy/gunicorn.service
[Unit]
Description=Gunicorn daemon for Azmooneh Django App
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/azmooneh
Environment="PATH=/var/www/azmooneh/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=backend.settings_production"
ExecStart=/var/www/azmooneh/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/gunicorn/access.log \
    --error-logfile /var/log/gunicorn/error.log \
    backend.wsgi:application

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
